# crontab
*/5 * * * * python3 /opt/script/notify.py
